# station/etc/systemd/system/nodm.service.d/local.conf
# See https://serverfault.com/a/803794
[Unit]
After=pulseaudio.service docker.service
Requires=pulseaudio.service docker.service
