# station/etc/systemd/system/nodm.service.d/require_network.conf
# See https://serverfault.com/a/803794
[Unit]
After=network.target network-online.target
Wants=network-online.target
Requires=network.target network-online.target
