# station/etc/systemd/system/nodm.service.d/station_with_projector.conf
# See https://serverfault.com/a/803794
[Unit]
After=projector-start.service
Requires=projector-start.service
