# station/etc/systemd/system/pulseaudio.service.d/restart.conf
# See https://serverfault.com/a/803794
[Service]
##Restart=on-failure
Restart=always
RestartSec=2
StartLimitInterval=300
StartLimitBurst=100


